# debian/db/crontab
30 3 * * * root        sudo -u postgres pg_dump civicboom | gzip --rsyncable > /root/civicboom.sql.gz
