# debian/db/crontab
30 3 * * * root sudo -u postgres pg_dump -Fc civicboom > /root/civicboom.pgdump
