# debian/api/postinst
#!/bin/sh

. /usr/share/debconf/confmodule

if [ ! -d /tmp/cb-website-temp ] ; then
	mkdir /tmp/cb-website-temp
	chown www-data:www-data /tmp/cb-website-temp
fi

touch /var/log/civicboom-pylons.log
chown www-data:www-data /var/log/civicboom-pylons.log

if [ "$1" = "configure" ] || [ "$1" = "reconfigure" ] ; then
	db_get cb-website-api/do_timed_tasks
	if [ "$RET" = "false" ] ; then
		rm -f /etc/cron.d/cb-website-api
	fi
fi

update-rc.d cb-website defaults
invoke-rc.d cb-website start

cd /opt/cb/share/website
find -name "*.pyc" | xargs rm -f

#echo "Copying static data to S3"
#/opt/cb/bin/cb-task sync_public_to_warehouse >/dev/null

db_stop
