# debian/crontab

30 3 * * * root        sudo -u postgres pg_dump civicboom | gzip --rsyncable > /root/civicboom.sql.gz
#30 2 * * * www-data    python /opt/cb/share/website/tasks.py
