# debian/postinst
#!/bin/sh

mkdir -p /var/log/apache2/civicboom.com/
make-ssl-cert generate-default-snakeoil

ln -sf /etc/nginx/sites-available/civicboom.com /etc/nginx/sites-enabled/

if [ ! -d /var/log/civicboom ] ; then
	mkdir /var/log/civicboom
	chown www-data:www-data /var/log/civicboom
fi

update-rc.d cb-website defaults
invoke-rc.d cb-website start
invoke-rc.d nginx restart

. /opt/cb-env/bin/activate
cd /opt/cb/share/website
echo "Copying static data to S3"
paster request production.ini /task/sync_public_to_warehouse >/dev/null
